# docker/kafka/supervisord.conf
[supervisord]
nodaemon=true

[program:kafka]
command=/usr/bin/kafka-start
autostart=true
autorestart=true